# capital_saas/deploy/daily_backup.service.example
[Unit]
Description=Capital SaaS Daily Backup

[Service]
WorkingDirectory=/opt/capital-saas/capital_saas
EnvironmentFile=/opt/capital-saas/.env
ExecStart=/opt/capital-saas/.venv/bin/python scripts/run_daily_backup.py
Type=oneshot
User=www-data
Group=www-data

